# asterisk-db-connector/etc/init/asterisk-db-connector.conf
start on [2345]
stop on [!2345]

exec python /opt/asterisk_db_connector/src/app.py
respawn
