# extra_install/asterisk_customs/asterisk-db-connector/DEBIAN/postinstall
#!/bin/bash

set -e

#systemctl enable asterisk_db_connector

echo "To start asterisk_db_connector, just run "
echo "service asterisk_db_connector start"

