# extra_install/asterisk_customs/asterisk-db-connector/etc/systemd/system/asterisk_db_connector.service
[Unit]
Description=Interface for database operations from mqtt

[Service]
Type=simple
ExecStart=python /opt/asterisk_db_connector/src/app.py

[Install]
WantedBy=multi-user.target
